# scripts/systemd/norman-tui-host-pressure-guard.timer
[Unit]
Description=Run Norman TUI host pressure guard every minute

[Timer]
OnBootSec=1m
OnUnitActiveSec=1m
AccuracySec=15s
Unit=norman-tui-host-pressure-guard.service

[Install]
WantedBy=timers.target
