# scripts/systemd/norman-agent-console-sync-local.timer
[Unit]
Description=Periodic local Norman agent console template sync

[Timer]
OnBootSec=2m
OnUnitInactiveSec=15m
Persistent=true
Unit=norman-agent-console-sync-local.service

[Install]
WantedBy=timers.target
