# scripts/systemd/norman-agent-console-sync.timer
[Unit]
Description=Periodic Norman agent console template sync

[Timer]
OnBootSec=2m
OnUnitInactiveSec=15m
Persistent=true
Unit=norman-agent-console-sync.service

[Install]
WantedBy=timers.target
