# scripts/systemd/norman-bbs-doctor.timer
[Unit]
Description=Periodic Norman BBS policy doctor

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true
Unit=norman-bbs-doctor.service

[Install]
WantedBy=timers.target
