# scripts/systemd/norman.service.d/10-runtime-env.conf
[Service]
EnvironmentFile=
EnvironmentFile=-/etc/norman/runtime-identities.env
