# scripts/systemd/norman.service.d/20-runtime-guardrails.conf
[Service]
LimitNOFILE=65536
ExecStart=
ExecStart=/home/kristopher/code/norman/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --limit-concurrency 192 --backlog 256 --timeout-keep-alive 5
